# 192.168.10.104/crontab
0   8 * * 2 /root/update_pkgs.sh
0   8 * * 5 /root/update_pkgs.sh
0   4 * * * /root/sync-vmail-backup.sh
